# shell编程学习/常用软件安装脚本/1.php脚本安装
#！bin/bash
echo "请输入你要安装的nginx版本(12,14,16,17)："
read ngver
echo "请输入你要安装的mysql版本(55,56,80)："
read myver
echo "请输入你要安装的php版本（55,56,70，72）:"
read phpver

#安装nginx
if [[ "$ngver" -eq 12 ]]
 then
	wget http://nginx.org/download/nginx-1.12.2.tar.gz
	tar -zxvf nginx-1.12.2.tar.gz
    cd nginx-1.12.2 
elif [[ "$ngver" -eq 14 ]]
 then
	wget http://nginx.org/download/nginx-1.14.2.tar.gz
	tar -zxvf nginx-1.14.2.tar.gz
    cd nginx-1.14.2 
elif [[ "$ngver" -eq 16 ]]
 then
	wget http://nginx.org/download/nginx-1.16.1.tar.gz
		tar -zxvf nginx-1.16.1.tar.gz
    cd nginx-1.16.1 
elif [[ "$ngver" -eq 17 ]]
 then
	wget http://nginx.org/download/nginx-1.17.3.tar.gz
	tar -zxvf nginx-1.17.3.tar.gz
    cd nginx-1.17.3
else
	echo "没有相应的版本，请重新运行脚本"
fi
yum install gcc gcc-c++ pcre pcre-devel zlib zlib-devel openssl openssl-devel -y
groupadd www
useradd -g www www
./configure --prefix=/usr/local/nginx --without-http_memcached_module --user=www --group=www --with-http_stub_status_module --with-http_ssl_module --with-http_gzip_static_module
make && make install 
/usr/local/nginx/sbin/nginx
echo '#!/bin/bash
# chkconfig: - 30 21
# description: http service.
# Source Function Library
. /etc/init.d/functions
# Nginx Settings
NGINX_SBIN="/usr/local/nginx/sbin/nginx"
NGINX_CONF="/usr/local/nginx/conf/nginx.conf"
NGINX_PID="/usr/local/nginx/logs/nginx.pid"
RETVAL=0
prog="Nginx"
start() {
echo -n $"Starting $prog: "
mkdir -p /dev/shm/nginx_temp
daemon $NGINX_SBIN -c $NGINX_CONF
RETVAL=$?
echo
return $RETVAL
}
stop() {
echo -n $"Stopping $prog: "
killproc -p $NGINX_PID $NGINX_SBIN -TERM
rm -rf /dev/shm/nginx_temp
RETVAL=$?
echo
return $RETVAL
}
reload(){
echo -n $"Reloading $prog: "
killproc -p $NGINX_PID $NGINX_SBIN -HUP
RETVAL=$?
echo
return $RETVAL
}
restart(){
stop
start
}
configtest(){
$NGINX_SBIN -c $NGINX_CONF -t
return 0
}
case "$1" in
start)
start
;;
stop)
stop
;;
reload)
reload
;;
restart)
restart
;;
configtest)
configtest
;;
*)
echo $"Usage: $0 {start|stop|reload| 
restart|configtest}"
RETVAL=1
esac
exit $RETVAL' > /etc/init.d/nginx

chkconfig --add /etc/init.d/nginx
chmod 755 /etc/init.d/nginx
chkconfig --add nginx
chkconfig nginx on
cd

#安装mysql
if [[ "$myver" -eq 55 ]]; then
	wget https://dev.mysql.com/get/Downloads/MySQL-5.5/MySQL-5.5.62-1.el6.src.rpm
elif [[ "$myver" -eq 56 ]]; then
	wget https://dev.mysql.com/get/Downloads/MySQL-5.6/MySQL-5.6.45-1.el6.src.rpm
elif [[ "$myver" -eq 57 ]]; then
	wget https://dev.mysql.com/get/Downloads/MySQL-5.7/mysql-community-5.7.27-1.el6.src.rpm
elif [[ "$myver" -eq 80 ]]; then
	wget https://dev.mysql.com/get/Downloads/MySQL-8.0/mysql-community-8.0.17-1.el6.src.rpm
else
    echo "没有相应的版本，请重新运行脚本" 
fi



#安装php
if [[ "$phpver" -eq 55 ]]; then
    wget https://www.php.net/distributions/php-5.5.30.tar.gz
elif [[ "$phpver" -eq 56 ]]; then
	wget https://www.php.net/distributions/php-5.6.40.tar.gz
elif [[ "$phpver" -eq 70 ]]; then
	wget https://www.php.net/distributions/php-7.0.26.tar.gz
elif [[ "$phpver" -eq 72 ]]; then
	wget https://www.php.net/distributions/php-7.2.20.tar.gz
else
	echo "没有相应的版本，请重新运行脚本"
fi
tar -zxvf php* && cd php*
yum -y install gcc gcc-c++ autoconf libjpeg libjpeg-devel libpng libpng-devel freetype freetype-devel libxml2 libxml2-devel glibc glibc-devel glib2 glib2-devel openssl openssl-devel openldap openldap-devel pcre pcre-devel
